# container-files/etc/supervisor.d/zabbix.conf
[program:zabbix-agent]
command = /config/ds.sh /run/zabbix/zabbix_agentd.pid zabbix_agentd -c /etc/zabbix_agentd.conf
autorestart = true